# release_GPL/firmware/apps/etc/inittab

::sysinit:/etc/rcS
#::respawn:/bin/login
::respawn:/bin/sh
::restart:/sbin/init
